# docker/s6/docker/run
#!/bin/sh
set -eu

exec /usr/local/bin/dockerd-run.sh
